# deployment/social-api-celery-worker.service
[Unit]
Description=Social Media Manager Celery Worker
After=network.target redis.service postgresql.service
Wants=redis.service postgresql.service

[Service]
Type=forking
User=root
Group=root
WorkingDirectory=/opt/social-media/backend
Environment=PATH=/opt/social-media/backend/venv/bin
EnvironmentFile=/opt/social-media/backend/.env
ExecStart=/opt/social-media/backend/venv/bin/celery -A social_backend worker --loglevel=info --pidfile=/run/celery/worker.pid --logfile=/var/log/celery-worker.log --detach
ExecStop=/bin/kill -s TERM $MAINPID
ExecReload=/bin/kill -s HUP $MAINPID
PIDFile=/run/celery/worker.pid
Restart=always
RestartSec=10

# Create PID directory
RuntimeDirectory=celery
RuntimeDirectoryMode=755

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/social-media /var/log /run/celery /tmp

[Install]
WantedBy=multi-user.target